# device/com.amit.marginalia-publish.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.amit.marginalia-publish</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/amit/Projects/marginalia/device/publish.sh</string>
    </array>
    <key>StartInterval</key>
    <integer>7200</integer>
    <key>RunAtLoad</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/amit/Projects/marginalia/.cache/publish.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/amit/Projects/marginalia/.cache/publish.stderr.log</string>
</dict>
</plist>
